# templates/com.tinglevpn.wg.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tinglevpn.wg</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/usr/local/etc/wireguard/wg-daemon.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <!-- CRITICO: sem esta chave o launchd mata o wireguard-go que o wg-quick
         lanca em background (daemonizado), ~70ms apos subir. Era a causa raiz
         do crash-loop. Com ela, o wireguard-go sobrevive como processo
         abandonado, igual a subir "wg-quick up wg0" no terminal. -->
    <key>AbandonProcessGroup</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/tinglevpn-wg.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/tinglevpn-wg.log</string>
</dict>
</plist>
